# files-to-deploy/additional-config-files/other/systemd-units/glances-web.service
[Unit]
Description=Glances

[Service]
ExecStart=/bin/glances --webserver
Restart=on-abort

[Install]
WantedBy=multi-user.target
