# files-to-deploy/additional-config-files/other/systemd-units/pia-tools-check.service
[Unit]
Description=VPN connection check
[Service]
Type=oneshot
ExecStart=/usr/bin/pia-tools --check
